# scripts/deb-scripts/postinst
#!/bin/bash
set -e

systemctl daemon-reload
systemctl enable byocvpn-daemon
systemctl start byocvpn-daemon

exit 0
