# scripts/com.byocvpn.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.byocvpn.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>mkdir -p /var/run/byocvpn/release &amp;&amp; chmod 755 /var/run/byocvpn/release &amp;&amp; exec /Library/PrivilegedHelperTools/byocvpn-daemon</string>
    </array>

    <!-- Start immediately when loaded, and on every boot -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Redirect logs to a known location -->
    <key>StandardOutPath</key>
    <string>/Library/Logs/byocvpn-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/Library/Logs/byocvpn-daemon.error.log</string>
</dict>
</plist>
